# dot_config/systemd/user/sunshine.service.d/override.conf.tmpl
{{- if eq .chezmoi.os "linux" -}}
[Service]
ExecStart=
ExecStart=/usr/local/bin/sunshine
{{- end -}}
